# linux/crontab/compress.sh
#!/bin/bash

tar -zcvf /home/ec2-user/test.tar.gz /home/ec2-user/test
